# go-judge.conf
[program:go-judge]
command=/opt/go-judge
autostart=true
autorestart=true
stderr_logfile=/var/log/go-judge.err.log
stdout_logfile=/var/log/go-judge.out.log